# tools/supervisor.conf
[program:goshort]
command=/usr/bin/env /usr/local/bin/goshort
autostart=true
autorestart=true
stderr_logfile=/var/log/goshort.err.log
stdout_logfile=/var/log/goshort.out.log